# scripts/nightly-live-smoke/substrate-smoke.timer
[Unit]
Description=Run the substrate nightly live smoke at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
